# ROOTFS/var/lib/dpkg/info/man-db.list
/var
/var/cache
/var/cache/man
/var/lib
/var/lib/man-db
/etc
/etc/manpath.config
/etc/cron.weekly
/etc/cron.weekly/man-db
/etc/cron.daily
/etc/cron.daily/man-db
