# ROOTFS/var/lib/dpkg/info/udev.list
/etc
/etc/udev
/etc/udev/udev.conf
/etc/modprobe.d
/etc/modprobe.d/fbdev-blacklist.conf
/etc/init.d
/etc/init.d/udev
/etc/init
/etc/init/udev.conf
/etc/init/udevmonitor.conf
/etc/init/udevtrigger.conf
